# docker/systemd/calibre-import.service
[Unit]
Description=Import downloaded books into Calibre library
After=docker.service

[Service]
Type=oneshot
# Replace <YOUR_USER> with the local user that owns the library dirs and is in the docker group
User=<YOUR_USER>
SupplementaryGroups=docker
ExecStart=/library/docker/calibre-import.sh
StandardOutput=journal
StandardError=journal
