# docker/systemd/calibre-status-server.service
[Unit]
Description=Calibre import status HTTP server
After=network.target

[Service]
Type=simple
User=nobody
Group=nogroup
# Replace <SERVER_IP> with the LAN IP of the host running this service
ExecStart=/usr/bin/python3 -m http.server 8999 --directory /var/lib/calibre-import --bind <SERVER_IP>
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
